# supervisord/dnsmasq.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq -k
autostart=true
autorestart=true
user=root
