# supervisord/haproxy.conf
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg -n 4096 -N 4096
autostart=true
autorestart=true
user=root
